# scripts/systemd/voiceawareness-ca.service
[Unit]
Description=voiceawareness.ca Node.js application
After=network.target

[Service]
Type=simple
User=voiceawarenessca
Group=psacln
WorkingDirectory=/var/www/vhosts/voiceawareness.ca/httpdocs
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=PATH=/opt/plesk/node/24/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
EnvironmentFile=-/var/www/vhosts/voiceawareness.ca/httpdocs/.env
ExecStart=/opt/plesk/node/24/bin/node app.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
